# dot_config/systemd/user/garmin-calendar.service
[Unit]
Description=Pull Garmin Calendar

[Service]
Type=simple
ExecStart=/home/m2/.local/bin/garmin-calendar-pull

[Install]
WantedBy=default.target
