# dot_config/systemd/user/gdrivebackup.service
[Unit]
Description=Pull Garmin Calendar

[Service]
Type=simple
ExecStart=/home/m2/.local/bin/gdrivebackup

[Install]
WantedBy=default.target
