# dot_config/systemd/user/gdrivebackup.timer
[Unit]
Description=Pull Garmin Calendar Weekly Timer

[Timer]
OnCalendar=Daily
Persistent=true

[Install]
WantedBy=default.target
